# presto/src/etc/service/presto/run
#!/bin/sh

exec /opt/presto/bin/launcher run >>/var/log/presto_stdoe.log 2>&1
