# hive/src/etc/service/hive/run
#!/bin/sh

exec ./hive -service metastore
#sleep 10000d
