# hive/src/etc/service/hive/run
#!/bin/sh

initialize() {
 	sleep 1
 	if [ -f /tmp/dbinit.sql ]; then
 		/opt/hive/bin/hive -f /tmp/dbinit.sql
 		rm /tmp/dbinit.sql
 	fi
}

initialize &

# Ugly hack
sleep 10000d
